# templates/etc/profile.d/flume-ng.sh.erb
export FLUME_HOME=<%= flume_home %>-<%= flume_version %>
export PATH=$FLUME_HOME/bin:$PATH