# files/etc/profile.d/flume-ng.sh
export FLUME_HOME=/opt/flume-ng-1.2.0
export PATH=$FLUME_HOME/bin:$PATH